# profiles/templates/3.5/3_ac_install_live/1-merge/sys-libs/glibc/client/cron.hourly/refresh-nscd
# Calculate chmod=0755 module(client)!=&&client.os_remote_auth!=
#!/bin/bash

# Updating of domain users nss-cache.
# The NSS-cache is used for fix resume from suspend mode of domain users
#
# Обновление кэша пользователей домена.
# Кэширование используется для поддержки спящего режима доменных машин.

if pgrep -f /usr/sbin/nscd &>/dev/null; then
    # get hours uptime of nscd
    working_hours=$(( ($(date +%s) - $(stat -c%Y /proc/`pgrep -f /usr/sbin/nscd`)) / 3600 ))
    # check multiplicity of 3 and check domain server available
    if [[ "$(( $working_hours % #-client.cl_client_nscd_cache-# ))" == "0" ]] &&
        ldapsearch -D "#-server(unix.bind_dn)-#" -w "#-server(unix.bind_pw)-#" -b "#-server(ldap.services_dn)-#"  -h "#-client.cl_remote_host-#" -l5 -onettimeout=5 -n &>/dev/null;then
        /usr/sbin/nscd-refresh
    fi
fi
