# profiles/templates/deprecated/3.1/1_ac_install_disk/2-system/sys-apps/sysvinit/inittab
# Calculate format=patch
<reg>^#(s[01])</reg>
<text>\1</text>
<reg>(c[0-9]:.*agetty).*(tty\d linux)</reg>
<text>\1 38400 \2</text>
