# profiles/templates/3.3/5_ac_client_undomain/sys-apps/sysvinit/inittab
# Calculate format=patch
<reg>l06:06:wait:/usr/sbin/logout
</reg>
<text></text>
