# profiles/templates/3.1/4_ac_install_live/networkmanager/dispatcher.d/ip-down.sh
# Calculate chmod=0755
#!/bin/bash

DEV_NAME="$1"
NM_ACTION="$2"

if [[ ${NM_ACTION} == "down" ]]
then
	:
fi
