# profiles/templates/3.6/6_ac_desktop_profile/2-user/lxqt-base/LXQt-0.10/.config/autostart/setdpi.desktop
[Desktop Entry]
Exec=#-ur_home_path-#/.config/autostart/setdpi.sh
Name=Set DPI
OnlyShowIn=LXQt;
Type=Application
Version=1.0
