# profiles/templates/3.3/3_ac_install_live/1-merge/sys-apps/openrc/local.d/fixsystemtime.start
# Calculate os_install_root_type==livecd&&pkg(sys-apps/openrc)>=0.8.0 comment=# chmod=0755
# fix system time on livecd
#?os_install_clock_type==local#
hwclock --hctosys --localtime
#os_install_clock_type#
#?os_install_clock_type!=local#
hwclock --hctosys --utc
#os_install_clock_type#
