# profiles/templates/deprecated/3.4/6_ac_desktop_profile/2-user/kde-base/KDE-4.10/.config/autostart/kres-migrator.desktop
# Calculate format=kde
[Desktop Entry]
Comment=Resources migrator
Exec=/usr/bin/kres-migrator
GenericName=Resources migrator
Icon=accessories-text-editor
Name=Resources migrator
