# profiles/templates/deprecated/3.4/5_ac_client_domain/sys-libs/glibc/local.d/nscd-restart.start
# Calculate chmod=0755
#!/bin/bash

pidfile="/var/run/nscd/nscd.pid"
if [ -f ${pidfile} ]
then
	start-stop-daemon --stop --quiet \
		--exec /usr/sbin/nscd --pidfile ${pidfile}

	start-stop-daemon --start --quiet --background \
		--exec /usr/sbin/nscd --pidfile ${pidfile} \
		-- $secure --config-file=/etc/nscd-cache.conf
fi
