# profiles/templates/3.6/1_ac_install_disk/2-system/00-calculate-install/update-pkg/update-pkg
# Calculate chmod=755 path=/etc/init.d cl_install_update_pkg_set==on
#!/sbin/openrc-run
# Copyright 2018 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="Update package at first boot"
description="Update packages at first boot"

depend() {
	want local
}

start() {
	ebegin "Update packages"
	cl-update -f &>>/dev/tty12
	res=$?
	LANG=en_US.utf8 rc-update del update-pkg default
	LANG=en_US.utf8 /sbin/rc-update --update
	eend $res
}
