# profiles/templates/deprecated/3.1/5_ac_client_undomain/sys-apps/sysvinit/inittab
# Calculate format=patch pkg(sys-apps/calculate-client)>=3.1.5-r4
<reg>l06:06:wait:/usr/sbin/logout
</reg>
<text></text>
