# profiles/templates/3.4/5_ac_client_domain/sys-apps/sysvinit/inittab
# Calculate format=patch
<reg>(rc::bootwait:/sbin/(?:open)?rc boot

)(l0:0:wait:/sbin/(?:open)?rc shutdown)</reg>
<text>\1l06:06:wait:/usr/sbin/logout
\2</text>
