# profiles/templates/3.1/5_ac_client_domain/sys-libs/glibc/cron.hourly/refresh-nscd
# Calculate chmod=0755
#!/bin/bash

# Updating of domain users nss-cache.
# The NSS-cache is used for fix resume from suspend mode of domain users
#
# Обновление кэша пользователей домена.
# Кэширование используется для поддержки спящего режима доменных машин.

if pgrep nscd &>/dev/null; then
    # get hours uptime of nscd
    working_hours=$(( ($(date +%s) - $(stat -c%Y /proc/`pgrep nscd`)) / 3600 ))
    # check multiplicity of 3 and check domain server available
    if [[ "$(( $working_hours % #-cl_client_nscd_cache-# ))" == "0" ]] &&
        ldapsearch -D "#-server(unix.bind_dn)-#" -w "#-server(unix.bind_pw)-#" -b "#-server(ldap.services_dn)-#"  -h "#-cl_remote_host-#" -l5 -onettimeout=5 -n &>/dev/null
        /usr/sbin/nscd-refresh
    fi
fi
