# profiles/templates/3.1/5_ac_client_domain/sys-libs/glibc/cron.hourly/refresh-nscd
# Calculate chmod=0755
#!/bin/bash

# Updating of domain users nss-cache.
# The NSS-cache is used for fix resume from suspend mode of domain users
#
# Обновление кэша пользователей домена.
# Кэширование используется для поддержки спящего режима доменных машин.

if pgrep nscd &>/dev/null; then
    # get hours uptime of nscd
    working_hours=$(( ($(date +%s) - $(stat -c%Y /proc/`pgrep nscd`)) / 3600 ))
    # check multiplicity of 3
    if [[ "$(( $working_hours % 3 ))" == "0" ]]; then
        /usr/sbin/nscd-refresh
    fi
fi
