# profiles/templates/install/1live/openrc/local.d/fixsystemtime.start
# Calculate os_install_root_type==livecd&&pkg(sys-apps/openrc)>=0.8.0 comment=# chmod=0755
# fix system time on livecd
hwclock --hctosys
