# profiles/templates/3.6/3_ac_install_live/2-laptop/lenovo/sys-apps/calculate-utils/init.d/laptopd
# Calculate comment=# chmod=0755
#!/sbin/runscript
# Copyright 1999-2010 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-power/acpid/files/acpid-2.0.3-init.d,v 1.1 2010/04/04 16:39:09 ssuominen Exp $

depend() {
	after *
}

start() {
	ebegin "Configure laptop"
	setkeycodes e007 187   # XF86Launch8
	setkeycodes e00e 185   # XF86Launch5
	setkeycodes e04e 0     # Drop plus-minux
	if lsmod | grep -q "^psmouse"
	then
		rmmod psmouse
		modprobe psmouse
	fi
	lsmod | grep -q "^acer_wmi" && rmmod acer-wmi
#?hr_laptop_model==Mariana-3B#
	lsmod | grep -q "^tg3" && rmmod tg3
	modprobe -q broadcom
	modprobe -q tg3
#hr_laptop_model#
	eend 0
}

stop() {
	eend 0
}

