# profiles/templates/3_ac_install_live/1-merge/sys-apps/sysvinit/inittab
# Calculate comment=# format=patch
#?module(client)!=&&client.os_remote_auth!=#
<reg>(rc::bootwait:/sbin/(?:open)?rc boot

|l0u:0:wait:/sbin/telinit u
)(l0:0:wait:/sbin/(?:open)?rc shutdown)</reg>
<text>\1l06:06:wait:/usr/sbin/logout
\2</text>
#!module#
<reg>l06:06:wait:/usr/sbin/logout
</reg>
<text></text>
#module#
