# profiles/templates/3.5/2_ac_install_merge/gnome-base/gvfs/network.mount
# Calculate format=kde path=/usr/share/gvfs/mounts
[Mount]
AutoMount=false
