# net-misc/haspd/files/hasplmd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend () {
    need net
    after hasplmd
}

start() {
	ebegin "Starting HASP SRM"
	eval start-stop-daemon --start --quiet --exec /opt/sbin/hasplmd -- ${HASPLMD_ARGS}
	eend $?
}

stop() {
	ebegin "Stopping HASP SRM"
	start-stop-daemon --stop --quiet --exec /opt/sbin/hasplmd
	eend $?
}
