# net-misc/haspd/files/skeyd
#!/sbin/runscript
# Copyright 2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="/opt/sbin/skeyd"
command_args="${SKEYD_ARGS}"
description="Smartkey 3 USB/LPT driver"

if [ "${RC_VERSION:-0}" = "0" ]; then
        start() {
		eerror "This script cannot be used for baselayout-1."
		return 1
	}
fi

depend () {
    need net aksusbd
}
