# net-misc/haspd/files/skeyd
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend () {
    need net
    after aksusbd
}

start() {
	ebegin "Starting Smartkey 3 USB/LPT driver"
	eval start-stop-daemon --start --quiet --exec /opt/sbin/skeyd
	eend $?
}

stop() {
	ebegin "Starting Smartkey 3 USB/LPT driver"
	eval start-stop-daemon --stop --quiet --exec /opt/sbin/skeyd
	eend $?
}
