# profiles/templates/3.3/6_ac_desktop_profile/2-user/kde-apps/KDE-5/.config/autostart/kres-migrator.desktop
# Calculate format=kde ur_domain_set==on
[Desktop Entry]
Comment=Resources migrator
Exec=/usr/bin/kres-migrator
GenericName=Resources migrator
Icon=accessories-text-editor
Name=Resources migrator
