# profiles/templates/deprecated/3.5/6_ac_backup_create/1-root/etc/NetworkManager/dispatcher.d/.calculate_directory
# Calculate append=join link=/etc/NetworkManager/dispatcher.d
